# api/scripts/systemd/divergram-db-backup.timer
[Unit]
Description=Run Divergram PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Seoul
Persistent=true
RandomizedDelaySec=10m
Unit=divergram-db-backup.service

[Install]
WantedBy=timers.target
