# api/scripts/systemd/divergram-health-monitor.service
[Unit]
Description=Divergram API health monitor
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/divergram/health-monitor.env
ExecStart=/usr/local/sbin/divergram-health-monitor
User=root
Group=root
